# net/bind-latest/files/hotplug.d/iface/19-bind
NAME=bind
COMMAND=/usr/sbin/named

[ "$ACTION" = "ifup" -a "$INTERFACE" = "wan" ] && {
        [ -x $COMMAND ] && [ -r $CONFIG ] && {
		killall named
		/etc/init.d/named start
        } &
}
