# net/bind-latest/files/xinetd.d/named
service domain
{
	disable         = no
	socket_type     = dgram
	wait            = yes
	user            = root
	group           = bind
	only_from	= 0.0.0.0/0 ::/0
	server          = /usr/sbin/named
	server_args     = -c /etc/bind/conf/named.conf -u bind -t /etc/chroot/named
	log_on_failure  += USERID
	flags		= IPv6
}


