# utils/nanom5poe/files/etc/init.d/poe
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 OpenWrt.org

START=01

boot() {
	poe start
}

start() {
	poe start
}

stop() {
	poe stop
}
