# net/aqm-scripts/files/etc/hotplug.d/iface/11-aqm
#!/bin/sh
[ "$ACTION" = ifup ] && /etc/init.d/aqm enabled && /usr/lib/aqm/generate.sh interface "$INTERFACE" | sh
