# net/bind-latest/files/named.prep
#!/bin/sh /etc/rc.common
# Setup the jail for bind

START=46
STOP=46

CHRT=/etc/chroot/named/
CHRTB=$CHRT/etc/bind

config_file=$CHRT/etc/bind/conf/named.conf
pid_file=$CHRT/var/run/named/named.pid

KEYFILES="us them rndc"
DEVFILES="/dev/log /dev/zero /dev/urandom /dev/random /dev/null"

gen_rndc() {
    local keyname=$1
    [ ! -s $CHRTB/keys/$keyname.key ] && {
    	rndc-confgen -u root -k $keyname-key -c /tmp/$keyname.key -r /dev/urandom -a
    	chgrp bind /tmp/$keyname.key
    	chmod og-rwx /tmp/$keyname.key
    	chmod g+r /tmp/$keyname.key
    	mv /tmp/$keyname.key $CHRTB/keys/$keyname.key
    	}
}

# When we start on squashfs, we're going to have bad perms, so fix them
# FIXME It would be better to do this without pounding flash, too

# perms generated via find . -printf '%h/%f %#m %u %g\n'

gen_perms() {

mkdir $CHRTB/managed $CHRTB/keys

cat /etc/bind.perms | while read fname mode owner group
do
    chgrp $group /etc/$fname
    chown $owner /etc/$fname
    chmod $mode /etc/$fname
done
}

setup_jail() {

    mkdir -p $CHRT/dev $CHRT/tmp 
#    mount -t tmpfs -o size=64k $CHRT/tmp 
#    mount -t tmpfs -o size=64k $CHRT/dev

    cp -a $DEVFILES $CHRT/dev

# if the keyfiles don't exist or are zero length, create them
  
    gen_rndc us
    gen_rndc them
    gen_rndc rndc
    cp $CHRTB/keys/rndc.key /etc
    gen_perms
}

start() {
	setup_jail
}

stop() {
killall named
# umount $CHRT/dev
}
